# supervisord.conf
[supervisord]
nodaemon=true

[program:php]
command=/usr/sbin/php5-fpm -F

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_events_enabled=true
stderr_events_enabled=true
